# homedir/.crontab
# brew nightly updates
0 3 * * *  brew update
0 4 * * *  brew upgrade
